# Framework2/saf/rootfs/etc/init.d/appmgr
#!/bin/sh /etc/rc.common
# Copyright (C) 2013-2014 OpenWrt.org

START=15
USE_PROCD=1
PROG=/sbin/appmgr

start_service(){
	procd_open_instance
	procd_set_param command "$PROG"
	procd_set_param respawn 3600 5 0
	procd_close_instance
}
