# Framework2/saf/rootfs/usr/lib/opkg/info/cloudclient.list
/etc/init.d/cloudclt
/etc/config/cloudclient
/etc/config/cloudclient.def
/sbin/cloudclient
